# examples/supervisor.conf
[program:citadel.reminder]
command=/opt/citadel.reminder/venv/bin/python3 /opt/citadel.reminder/main.py
autorestart=true